# files/rc.local
#!/bin/bash

if [ -e /opt/harbor/setup/.ran_customization ]; then
    exit
else
    /opt/harbor/setup/setup.sh &> /var/log/bootstrap.log
fi